# supervisor/etc/supervisor.d/nginx.conf
[program:nginx]
command=nginx
autostart=true
autorestart=true
startsecs=1
