# supervisor/etc/supervisor.d/php-fpm.conf
[program:php-fpm]
command=/etc/init.d/php-fpm start
autostart=true
autorestart=true
startsecs=1